# root/etc/services.d/freepbx-ldap/run
#!/usr/bin/with-contenv bash

exec freepbx-ldap 2>&1
